# push.cron
MAILTO=""
0,30 * * * * /home/linktrai/thomaspush.sh >& /dev/null
